# scripts/hermes-gateway.service
[Unit]
Description=Hermes Agent gateway (Nous Research)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=hermes
Group=hermes
Environment=HOME=/srv/hermes
Environment=XDG_RUNTIME_DIR=/run/user/999
Environment=DOCKER_HOST=unix:///run/user/999/docker.sock
Environment=PATH=/srv/hermes/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStartPre=/bin/sh -c 'uid=$(id -u hermes); i=0; until test -S "/run/user/$uid/docker.sock"; do i=$((i+1)); [ $i -ge 30 ] && exit 1; sleep 1; done'
ExecStart=/srv/hermes/.local/bin/hermes gateway
Restart=on-failure
RestartSec=5
TimeoutStopSec=210

NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/srv/hermes /var/lib/hermes /var/log/hermes
PrivateTmp=yes
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictSUIDSGID=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
CapabilityBoundingSet=
AmbientCapabilities=
SystemCallFilter=@system-service
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
